# .misc/Fedora/autostart/.config/autostart/sxhkd.desktop
[Desktop Entry]
Exec=sxhkd
Icon=
Name=sxhkd
Path=
Terminal=False
Type=Application
