# packages/autostart/chrome-cifhbcnohmdccbgoicgdjpfamggdegmo-Profile_1.desktop
[Desktop Entry]
Exec=/opt/google/chrome/google-chrome "--profile-directory=Profile 1" --app-id=cifhbcnohmdccbgoicgdjpfamggdegmo --force-dark-mode --enable-features=WebUIDarkMode --force-device-scale-factor=1.125
Icon=chrome-cifhbcnohmdccbgoicgdjpfamggdegmo-Profile_1
Name=Microsoft Teams
StartupWMClass=Teams
Terminal=false
Type=Application
Version=1.0
