# packages/autostart/keynav.desktop
[Desktop Entry]
Exec=keynav daemonize
Icon=
Name=keynav
Path=
Terminal=False
Type=Application
