# Fedora/autostart/.config/autostart/restart-audio.desktop
[Desktop Entry]
Exec=restart-audio
Icon=
Name=restart-audio
Path=
Terminal=False
Type=Application
